# docker-container/etc/profile.d/100-bme-activation.sh
# Global custom path
export PATH=/opt/bme:$PATH

# virtualenvwrapper activation
if [ -r /usr/share/virtualenvwrapper/virtualenvwrapper.sh ]; then
	source /usr/share/virtualenvwrapper/virtualenvwrapper.sh
fi
